# root/etc/s6-overlay/s6-rc.d/init-eden-config/up
#!/bin/execlineb -P
foreground { redirfd -w 1 /proc/1/fd/1 redirfd -w 2 /proc/1/fd/2 fdmove -c 2 1 /etc/s6-overlay/s6-rc.d/init-eden-config/init.sh }
